# systemd/.config/systemd/user/pa-applet.service
[Unit]
Description=PulseAudio System Tray
[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/pa-applet
[Install]
WantedBy=multi-user.target
